# config/webcore.service.erb
[Unit]
Description=Webcore Server
After=postgresql.service

[Service]
WorkingDirectory=<%= release_path %>
EnvironmentFile=/etc/www/webcore.env
ExecStart=<%= shared_path %>/webcore.sh
Type=simple
User=www
Group=www
Restart=always

[Install]
WantedBy=multi-user.target